# deploy/ai-monitor.service
[Unit]
Description=AI Monitor - Social Media & Repository Watch Platform
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=monitor
Group=monitor
WorkingDirectory=/opt/monitor

# Python environment
Environment=PYTHONPATH=/opt/monitor
EnvironmentFile=/opt/monitor/.env

# ExecStart - uses uvicorn on localhost only (nginx reverse proxy)
ExecStart=/opt/monitor/.venv/bin/python -m uvicorn ai_monitor.main:app \
    --host 127.0.0.1 --port 8081 \
    --workers 2 \
    --limit-concurrency 20 \
    --timeout-keep-alive 30

# Restart policy
Restart=always
RestartSec=5

# Logging
StandardOutput=append:/var/log/ai-monitor.log
StandardError=append:/var/log/ai-monitor-error.log

# Security
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
